# extra/oxidized.upstart
# by @andrewpwade in issue #45

start on started networking

respawn

setuid oxidized
setgid oxidized

chdir /home/oxidized

env HOME=/home/oxidized

pre-start script
    test -x /usr/local/bin/oxidized || { stop; exit 0; }
end script

exec /usr/local/bin/oxidized
